# packaging/linux/challenger-siem-agent.service
[Unit]
Description=Challenger SIEM Linux Agent
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=notify
User=challenger-siem
Group=challenger-siem
Environment=CHALLENGER_SIEM_AGENT_CONFIG=/etc/challenger-siem-agent/agentsettings.json
ExecStart=/opt/challenger-siem-agent/Challenger.Siem.LinuxAgent
Restart=on-failure
RestartSec=10s
UMask=0077
NoNewPrivileges=yes
CapabilityBoundingSet=
AmbientCapabilities=
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
ReadWritePaths=/var/lib/challenger-siem-agent
RuntimeDirectory=challenger-siem-agent
RuntimeDirectoryMode=0700
StateDirectory=challenger-siem-agent
StateDirectoryMode=0700
TasksMax=128
MemoryMax=250M

[Install]
WantedBy=multi-user.target
